# dist/autodeploy_technochat.service
[Unit]
Description=Auto-deploy technochat from master branch
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/technochat
ExecStart=/opt/technochat/dist/autodeploy_technochat.sh
EnvironmentFile=/etc/default/deploy_technochat
