# dist/autodeploy_technochat_rc.service
[Unit]
Description=Auto-deploy technochat RC
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/technochat
ExecStart=/opt/technochat/dist/autodeploy.sh --rc
EnvironmentFile=/etc/default/autodeploy_technochat
